# deploy/systemd/cm-backup.timer
[Unit]
Description=creative_manager DB バックアップを毎日

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
