# deploy/systemd/cm-search.service
[Unit]
Description=creative_manager 意味検索 (cm-search :8788・埋め込み)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/creative_manager/apps/worker
EnvironmentFile=-%h/.config/creative-manager.env
Environment=CM_DB=%h/projects/creative_manager/data/cm.sqlite
ExecStart=/usr/bin/env bash -lc 'exec uv run cm-search'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
