# deploy/systemd/cm-worker.service
[Unit]
Description=creative_manager ジョブワーカー (claude -p 消化・reap は api 側)
# api(DB作成・job表権威)が上がってから。agentic は creative-manager stdio を claude -p が自前 spawn。
After=cm-api.service
Wants=cm-api.service

[Service]
Type=simple
WorkingDirectory=%h/projects/creative_manager/apps/worker
EnvironmentFile=-%h/.config/creative-manager.env
Environment=CM_DB=%h/projects/creative_manager/data/cm.sqlite
Environment=CM_MCP_STDIO_CMD=pnpm
Environment=CM_MCP_STDIO_ARGS=["--filter","@cm/api","mcp"]
# worker は CM_DB を直接読む（job表は api と同一FK DDLでブートストラップ可）。HTTP api は不要。
ExecStart=/usr/bin/env bash -lc 'exec uv run cm-worker'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
